# apps/bash/profile.d/cargo.sh
# profile.d/cargo.sh -- Rust/Cargo.
# shellcheck shell=bash
path::append "$HOME/.cargo/bin"
