# apps/bash/profile.d/jetbrains.sh
# profile.d/jetbrains.sh -- JetBrains Toolbox shell scripts (idea, pycharm, ...).
# shellcheck shell=bash
export JETBRAINS="$HOME/.local/share/JetBrains/Toolbox/scripts"
path::append "$JETBRAINS"
